# supervisord.conf


[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

# ---------- FastAPI Backend (internal only) ----------
[program:fastapi]
directory=/app
command=uvicorn main:app --host 127.0.0.1 --port 8000 --workers 1
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

# ---------- Streamlit Frontend (public-facing) ----------
[program:streamlit]
directory=/app
command=streamlit run streamlit_app.py
    --server.port=%(ENV_PORT)s
    --server.address=0.0.0.0
    --server.headless=true
    --browser.gatherUsageStats=false
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
